# thingengine-server.service
# This file is part of ThingEngine
#
# Copyright 2015 Giovanni Campagna <gcampagn@cs.stanford.edu>
#
# See COPYING for details

[Unit]
Description=ThingEngine personal IoT service

[Service]
ExecStart=/opt/nodejs/bin/node /home/pi/thingengine-platform-server/main.js

[Install]
WantedBy=default.target
